# scripts/launchd/com.aidanny.calendar.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Calendar watcher — pre-call brief auto-trigger.
  Runs every 15 minutes, fires a pre-call brief for meetings starting in ≤30 min.

  IMPORTANT: After loading, grant Calendar access when macOS prompts.
  Or manually: System Settings → Privacy & Security → Calendars → enable Terminal.

  Install:
    cp scripts/launchd/com.aidanny.calendar.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.aidanny.calendar.plist

  Uninstall:
    launchctl unload ~/Library/LaunchAgents/com.aidanny.calendar.plist
    rm ~/Library/LaunchAgents/com.aidanny.calendar.plist

  Run now (test):
    launchctl start com.aidanny.calendar

  View logs:
    tail -f /tmp/cal-watcher.log

  Dry run (no API calls):
    node scripts/watch-calendar.mjs --dry-run
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.aidanny.calendar</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>/Users/danielpaul/Documents/Second Brain Obsidian App/scripts/watch-calendar.mjs</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/danielpaul/Documents/Second Brain Obsidian App</string>

  <!-- Every 15 minutes -->
  <key>StartInterval</key>
  <integer>900</integer>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/tmp/cal-watcher.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/cal-watcher.log</string>

  <key>KeepAlive</key>
  <false/>
</dict>
</plist>
